# daemons/systemd-pushup-api.service
# Example systemd --user unit for push up stats API backend
# Install as: ~/.config/systemd/user/pushup-stats-api.service

[Unit]
Description=Pushup Stats API
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.openclaw/workspace/pushup-stats-service
Environment=NODE_ENV=production
Environment=PORT=8788
ExecStart=/usr/bin/node dist/api/main.js
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
